# deploy/qdrant-chai.service
[Unit]
Description=Qdrant Vector Database for Chai
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker stop qdrant-chai
ExecStartPre=-/usr/bin/docker rm qdrant-chai
ExecStart=/usr/bin/docker run --rm --name qdrant-chai \
    -p 127.0.0.1:6333:6333 \
    -p 127.0.0.1:6334:6334 \
    -v /opt/qdrant_storage:/qdrant/storage \
    qdrant/qdrant
ExecStop=/usr/bin/docker stop qdrant-chai
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
